# gunicorn.conf
[program:gunicorn]
command=/usr/local/bin/gunicorn --timeout 600 main:app -b localhost:5000
directory=/app
